# images/debian-wheezy/patches/etc/init/sync-kernel-modules.conf
# sync kernel modules
description "synchronizes kernel modules"
author "Online Labs Cloud Team <cloud-team@labs.online.net>"

start on mounted MOUNTPOINT=/
kill timeout 300

task

script

/usr/sbin/oc-sync-kernel-modules
end script
